# deploy/hotel-bot.service
[Unit]
Description=Hotel Villa Margaritas WhatsApp bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/hotel-bot
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=-/opt/hotel-bot/.env
Environment=BOT_ID=principal
Environment=BOT_LABEL=Bot principal
Environment=BOT_AUTH_DIR=auth

[Install]
WantedBy=multi-user.target
